# deploy/yt-dlp-mcp-update.service
[Unit]
Description=yt-dlp-mcp: pull a fresh yt-dlp release
Documentation=https://github.com/wildcar/yt-dlp-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=movie
Group=movie
WorkingDirectory=/opt/yt-dlp-mcp
# `uv pip install -U yt-dlp` bumps the venv-pinned binary; the
# subsequent restart picks it up. Failures are non-fatal — the timer
# retries tomorrow, and the existing yt-dlp keeps running.
ExecStart=/usr/bin/env bash -c '/opt/yt-dlp-mcp/.venv/bin/python -m pip install -U --quiet yt-dlp'
ExecStartPost=/bin/systemctl restart yt-dlp-mcp.service
