# deploy/yt-dlp-mcp-update.timer
[Unit]
Description=Daily yt-dlp update for yt-dlp-mcp

[Timer]
# Once a day, with a randomized delay so multiple boxes on the same
# subnet don't all hammer pypi at the same minute.
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
